# udev/rules.d/44-usb-c-refresh.rules
ACTION=="add|remove", SUBSYSTEM=="usb", ATTRS{idVendor}=="001f", ATTRS{idProduct}=="0b21", RUN+="/usr/local/bin/keys-usb-handler.sh"
