# udev/services/wacom.service
[Unit]
Description=wacom configurations

[Service]
Type=oneshot
ExecStart=/home/omar/scripts/wacom.sh
