# udev/services/wacom.service
[Unit]
Description=Wacom configurations
After=graphical.target

[Service]
Type=oneshot
User=omar
ExecStart=/home/omar/scripts/wacom.sh

